# rc.local
#!/bin/bash

/usr/bin/uwsgi --emperor /home/bakovka4/django/bakovka4/vassals --daemonize /var/log/uwsgi/emperor.log --pidfile /tmp/uwsgi-emperor.pid

exit 0
